# systemd/pi-watchdog.service
[Unit]
Description=Process Watchdog Service for AdProcess

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 __PIWATCHDOG_SCRIPT__